# configuration/chat.conf
[program:chat]
command=python3 /root/chat/chat.py
autostart=true
autorestart=true
